# sshtunnel@.service
[Unit]
Description=SSH Tunnel managment
Requires=network.target
After=network.service

[Service]
Type=simple
ExecStart=/usr/bin/sshtunnel-systemd %i
